# .config/autostart/picom.desktop
[Desktop Entry]
Type=Application
Name=Picom
Comment=An X compositor
Exec=/usr/bin/picom
