# .config/autostart/xsettingsd.desktop
[Desktop Entry]
Type=Application
Name=xsettingsd
Comment=Simple Xsettings daemon
Exec=/usr/bin/xsettingsd
