# .config/autostart/xss-lock.desktop
[Desktop Entry]
Type=Application
Name=xss-lock
Comment=Lockscreen daemon
Exec=/usr/bin/xss-lock /usr/bin/sxlock
